# daemon/pkg/daemon/testdata/cron_corpus/cron.d/app
MAILTO=ops@example.com
TZ=America/New_York

15,45 * * * * app /opt/app/job --mode=batch
